# lightning-talks/2014/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx"

#[program:elasticsearch]
#command=/opt/elasticsearch/bin/elasticsearch
